# deploy/systemd/travel-wtatil-price-sync.timer
# @deprecated travel-wtatil-sync.timer kullanın. Bu timer aynı auto script'i çağırır.

[Unit]
Description=Travel Wtatil senkron (günlük, legacy unit adı)

[Timer]
OnBootSec=20min
OnCalendar=*-*-* 02:30:00
AccuracySec=15min
RandomizedDelaySec=45min
Persistent=true
Unit=travel-wtatil-price-sync.service

[Install]
WantedBy=timers.target
