# deploy/systemd/travel-social-worker.timer
# travel-social-worker.service ile birlikte kullanın.

[Unit]
Description=Travel sosyal medya worker zamanlayıcısı (10 dk)

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
AccuracySec=30s
RandomizedDelaySec=60s
Unit=travel-social-worker.service

[Install]
WantedBy=timers.target
