# deploy/systemd/travel-wtatil-sync.service
# Wtatil otomatik senkron — günlük (travel-wtatil-sync.timer).
#
# Kurulum:
#   sudo cp deploy/systemd/travel-wtatil-sync.service /etc/systemd/system/
#   sudo cp deploy/systemd/travel-wtatil-sync.timer /etc/systemd/system/
#   sudo chmod +x /var/www/vhosts/rezervasyonyap.tr/httpdocs/deploy/scripts/sync-wtatil-auto.sh
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now travel-wtatil-sync.timer
#
# Log: journalctl -u travel-wtatil-sync.service -n 200 --no-pager

[Unit]
Description=Travel — Wtatil tur otomatik senkron (dönem+fiyat)
After=network-online.target travel-api.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/vhosts/rezervasyonyap.tr/httpdocs
ExecStart=/bin/bash /var/www/vhosts/rezervasyonyap.tr/httpdocs/deploy/scripts/sync-wtatil-auto.sh
Environment=WTATIL_AUTO_IMPORT_NEW_LIMIT=9999
Nice=10
IOSchedulingClass=best-effort

[Install]
WantedBy=multi-user.target
