# deploy/systemd/travel-travelrobot-sync.service
[Unit]
Description=Travelrobot otomatik senkron (tur · otel incremental · uçuş)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/vhosts/rezervasyonyap.tr/httpdocs
Environment=TRAVEL_DB_ENV=/etc/rezervasyonyap/backend.env
# +x eksik olsa bile çalışsın (203/EXEC önlemi)
ExecStart=/bin/bash /var/www/vhosts/rezervasyonyap.tr/httpdocs/deploy/scripts/sync-travelrobot-auto.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
