# deploy/systemd/travel-ai-worker.timer
# travel-ai-worker.service ile birlikte kullanın (bkz. travel-ai-worker.service başlığı).
# Önceki: 3 dk → sürekli poll ekstra I/O; 10 dk yeterli.

[Unit]
Description=Travel AI worker zamanlayıcısı

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
AccuracySec=1min
RandomizedDelaySec=45s
Unit=travel-ai-worker.service

[Install]
WantedBy=timers.target
