# deploy/systemd/travel-yolcu360-sync.service
[Unit]
Description=Yolcu360 araç kiralama otomatik import
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/vhosts/rezervasyonyap.tr/httpdocs
Environment=TRAVEL_DB_ENV=/etc/rezervasyonyap/backend.env
# +x eksik olsa bile çalışsın (203/EXEC önlemi)
ExecStart=/bin/bash /var/www/vhosts/rezervasyonyap.tr/httpdocs/deploy/scripts/import-yolcu360-cars.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
