# deploy/systemd/travel-warm-cache.timer
# Vitrin ISR önbelleğini düzenli ısıtır (kritik liste/landing sayfaları).
# Boot/restart sonrası birkaç dakika içinde ilk ısıtma; sonra periyodik tazeleme.
# ISR revalidate (60-300s) zaten taze tutar; bu timer eviction/restart sonrası
# soğuk render penceresini kapatmak içindir.

[Unit]
Description=Travel — vitrin ISR önbellek ısıtma (5 dk)

[Timer]
OnBootSec=20s
OnUnitActiveSec=5min
AccuracySec=15s
Persistent=true
Unit=travel-warm-cache.service

[Install]
WantedBy=timers.target
