# deploy/systemd/travel-warm-cache.timer
# Vitrin ISR önbelleğini düzenli ısıtır (kritik liste/landing sayfaları).
# Önceki: 5 dk → katalog sorgusu + disk I/O yükü yüksekti.
# ISR revalidate (60-300s) zaten taze tutar; 45 dk yeter.

[Unit]
Description=Travel — vitrin ISR önbellek ısıtma (45 dk)

[Timer]
OnBootSec=45s
OnUnitActiveSec=45min
AccuracySec=1min
Persistent=true
RandomizedDelaySec=2min
Unit=travel-warm-cache.service

[Install]
WantedBy=timers.target
