# deploy/systemd/travel-web.service
# Travel frontend (Next.js). Kurulum: /etc/systemd/system/travel-web.service
#
# Sunucuda:
#   sudo cp deploy/systemd/travel-web.service /etc/systemd/system/travel-web.service
#   sudo mkdir -p /etc/rezervasyonyap
#   sudo cp deploy/systemd/frontend.env.example /etc/rezervasyonyap/frontend.env
#   sudo systemctl daemon-reload
#   sudo systemctl restart travel-web
#   sudo systemctl status travel-web
#
# --- Manuel deploy (WorkingDirectory ile AYNI klasorde build sart) ------------
# Asagidaki WorkingDirectory icinde ikinci kez "cd frontend" YAPMAYIN;
# .next baska yolda kalir -> Loading chunk failed / 500 / text/plain.
#
#   cd /var/www/vhosts/rezervasyonyap.tr/httpdocs/frontend
#   git pull origin main
#   rm -rf .next && npm run build
#   systemctl restart travel-web
#
# Git kokunuz httpdocs ise once oraya cd edin, sonra:
#   cd frontend && rm -rf .next && npm run build

[Unit]
Description=Travel Next.js (rezervasyonyap.tr)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/vhosts/rezervasyonyap.tr/httpdocs/frontend

# .env.local systemd tarafinda otomatik okunmaz; uretim env bu dosyadan gelir.
EnvironmentFile=/etc/rezervasyonyap/frontend.env
Environment=NODE_ENV=production
Environment=PORT=3000

# Next production server — server.mjs render-blocking CSS’i erteler (PSI).
# Kapatmak için EnvironmentFile içinde TRAVEL_DEFER_CSS=0
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/node ./server.mjs

Restart=always
RestartSec=5
KillMode=control-group
StandardOutput=append:/var/log/travel-web.log
StandardError=append:/var/log/travel-web.log

[Install]
WantedBy=multi-user.target
