# deploy/systemd/travel-vitrin-price-refresh.timer
# Vitrin fiyat önbelleğini günde 2 kez tazeler (yoğun saatlerde değil).
# Önceki: her 10 dk → 65k+ ilan full UPDATE + ağır disk I/O (datacenter uyarısı).
# Import sonrası hâlâ elle çalıştırılabilir: ./deploy/scripts/refresh-vitrin-prices.sh

[Unit]
Description=Travel — vitrin_price önbellek tazeleme (günde 2)

[Timer]
OnCalendar=*-*-* 01:15:00
OnCalendar=*-*-* 13:15:00
AccuracySec=5min
Persistent=true
RandomizedDelaySec=10min
Unit=travel-vitrin-price-refresh.service

[Install]
WantedBy=timers.target
