# deploy/systemd/travel-db-guard.timer
[Unit]
Description=Travel PostgreSQL connection guard timer

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=travel-db-guard.service

[Install]
WantedBy=timers.target
